# 7.1/alpine/clarin/supervisor/solr.conf
[program:solr]
command=solr-foreground
user=solr
autostart=true
autorestart=true
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
