# 7.1/alpine/supervisor/dcron.conf
[program:solr]
command=docker-entrypoint.sh solr-foreground
autostart=true
autorestart=true
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
